# redmine-backup.crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# m h dom mon dow 
0 2 * * * root /opt/redmine/scripts/backup.sh 2>&1 >> /var/log/redmine-backup.log
